# bin/Servers/Kylix/Daemon/bin/astaioserver
#!/bin/sh

# chkconfig: 345 20 80

# description: Start/Stop AstaIO Server
#
# This file belongs in /etc/init.d where it will be run
# on system startup and shutdown

# This init script contains functions specific for redhat
# and mandrake init scripts.

# Source function library.
# the . runs the script in the current process
. /etc/rc.d/init.d/functions

AstaIORunUser=mentcher

CURDATE=`date +%Y-%m-%d_%H_%M_%S`

export IB_LIB="/opt/interbase/lib"
export ASTAIO_SERVER_DIR="/home/mentcher/projects/astadaemon/bin"
export SERVER_CFG_FILE="$ASTAIO_SERVER_DIR/server.conf"
export SERVER_LOG_FILE="$ASTAIO_SERVER_DIR/server.log"

if [ -n "$LD_LIBRARY_PATH" ]; then
  export LD_LIBRARY_PATH="$ASTAIO_SERVER_DIR:$IB_LIB:$LD_LIBRARY_PATH"
else
  export LD_LIBRARY_PATH="$ASTAIO_SERVER_DIR:$IB_LIB"
fi

echo "LD_LIBRARY_PATH="$LD_LIBRARY_PATH
echo "PATH="$PATH

# Check the file is there and is executable.
# The [] represents the test command.  See man test.
[ -x $ASTAIO_SERVER_DIR/gdserver ] || exit 0

# See how we were called.
case "$1" in
  start)
	echo -n "Starting AstaIO server: "

  # move and zip the log file
  if [ -w $SERVER_LOG_FILE ]; then
    mv $SERVER_LOG_FILE $ASTAIO_SERVER_DIR/logs/serverlog_$CURDATE.log
    gzip $ASTAIO_SERVER_DIR/logs/serverlog_$CURDATE.log
  fi

  # the daemon function is in /etc/rc.d/init.d/functions
  # this will run the serer as the AstaIORunUser instead of root
	daemon --user $AstaIORunUser $ASTAIO_SERVER_DIR/gdserver
	RETVAL=$?
	;;
  stop)
	echo -n "Stopping AstaIO server: "

	killproc gdserver
	RETVAL=$?
	;;
  status)
	status gdserver
	RETVAL=$?
	;;
  restart|reload)
	$0 stop
	$0 start
	RETVAL=$?
	;;
  *)
	echo "Usage: astaioserver {start|stop|status|restart|reload}"
	exit 1
esac
echo

exit $RETVAL
